# flare.cron
0,30 * * * * mkdir -p ./logs/test_1-0-0 && docker exec flare-backend python3 /app/backend/flareRunner.py -c /app/data/cspec/test_1-0-0.json >> ./logs/test_1-0-0/$(date "+\%Y")_$(date "+\%m").log 2>> ./logs/CRON.log
*/15 * * * * mkdir -p ./logs/Laguna-Madre_Water-Level_Air-Temperature_120hrs && docker exec flare-backend python3 /app/backend/flareRunner.py -v  -c /app/data/cspec/Laguna-Madre_Water-Level_Air-Temperature_120hrs.json >> ./logs/Laguna-Madre_Water-Level_Air-Temperature_120hrs/$(date "+\%Y")_$(date "+\%m").log 2>> ./logs/error.log
*/15 * * * * mkdir -p ./logs/TWC-Laguna-Madre_Air-Temperature-Predictions_240hrs && docker exec flare-backend python3 /app/backend/flareRunner.py -v  -c /app/data/cspec/TWC-Laguna-Madre_Air-Temperature-Predictions_240hrs.json >> ./logs/TWC-Laguna-Madre_Air-Temperature-Predictions_240hrs/$(date "+\%Y")_$(date "+\%m").log 2>> ./logs/error.log
*/15 * * * * mkdir -p ./logs/MRE_Bird-Island_Water-Temperature_Ribbon && docker exec flare-backend python3 /app/backend/flareRunner.py -v  -c /app/data/cspec/MRE_Bird-Island_Water-Temperature_Ribbon.json >> ./logs/MRE_Bird-Island_Water-Temperature_Ribbon/$(date "+\%Y")_$(date "+\%m").log 2>> ./logs/error.log
*/15 * * * * mkdir -p ./logs/TWC-NDFD-Laguna-Madre_Air-Temperature-Predictions_240hrs && docker exec flare-backend python3 /app/backend/flareRunner.py -v  -c /app/data/cspec/TWC-NDFD-Laguna-Madre_Air-Temperature-Predictions_240hrs.json >> ./logs/TWC-NDFD-Laguna-Madre_Air-Temperature-Predictions_240hrs/$(date "+\%Y")_$(date "+\%m").log 2>> ./logs/error.log
# DO NOT DELETE THE EXTRA LINE BELOW!